# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:web]
user=root
command=bundle exec puma -C config/puma.rb
directory=/app
autorestart=true

[program:sidekiq]
user=root
command=bundle exec sidekiq
directory=/app
autorestart=true